# templates/flink.service.erb
[Unit]
Description=<%= scope['flink::service_name'] %>
After=network.target

[Service]
Type=forking
WorkingDirectory=<%= scope['flink::install_dir'] %>
ExecStart=<%= scope['flink::install_dir'] %>/bin/start-cluster.sh
ExecStop=<%= scope['flink::install_dir'] %>/bin/stop-cluster.sh
PIDFile=/var/run/flink-jobmanager.pid

[Install]
WantedBy=multi-user.target